# deploy/systemd/opencuttles-vision.service
[Unit]
Description=OpenCuttles vision sidecar (Florence-2 VLM for grounding + assertions)
After=network.target

[Service]
Type=simple
# Adjust User/WorkingDirectory to where the vision project is checked out.
User=atorres
WorkingDirectory=/home/atorres/opencuttles/vision
# HuggingFace model id. Swap to microsoft/Florence-2-large for more accuracy.
Environment=OPENCUTTLES_VISION_MODEL=microsoft/Florence-2-base
ExecStart=/home/atorres/opencuttles/vision/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8791
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
