# deploy/systemd/opencuttles-api.service
[Unit]
Description=OpenCuttles API
Documentation=https://github.com/opencuttles/opencuttles
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opencuttles
Group=opencuttles
SupplementaryGroups=kvm
WorkingDirectory=/var/lib/opencuttles
EnvironmentFile=-/etc/opencuttles/opencuttles.env
ExecStart=/opt/opencuttles/bin/opencuttles-api
Restart=on-failure
RestartSec=5s
# The API launches Cuttlefish AVDs (launch_cvd/crosvm) as child processes, so
# they live in this unit's cgroup. With the default KillMode=control-group, a
# stop/restart (e.g. a binary upgrade) would SIGKILL the whole cgroup and take
# every running AVD down with it. KillMode=process signals only the main API
# process on stop, leaving the AVDs running; the API reconnects to them over ADB
# on restart.
KillMode=process
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/opencuttles /var/log/opencuttles /run/opencuttles
RuntimeDirectory=opencuttles
RuntimeDirectoryMode=0750
UMask=0077

[Install]
WantedBy=multi-user.target
